# packaging/systemd/ydotoold.service
[Unit]
Description=ydotool daemon for input simulation
Documentation=man:ydotool(1)
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/ydotoold --socket-path=/run/ydotool/socket --socket-perm=0666
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
